# speakbox/deploy/speakbox.service
[Unit]
Description=Speakbox TTS Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/speakbox
EnvironmentFile=/opt/speakbox/.env
ExecStart=/opt/speakbox/speakbox
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
